# templates/seaweedfs-master.service.j2
[Unit]
Description=SeaweedFS master
After=network.target

[Service]
Type=simple
Restart=always
ExecStart={{ weed.location }}/weed master -ip.bind {{ weed.bind }} -ip {{ weed.ip }} \
                                          -mdir {{ weed.master.dir }} -defaultReplication {{ weed.defaultReplication }}
User={{ weed.user.name }}
SyslogIdentifier=weed-filer

[Install]
WantedBy=multi-user.target
